# scripts/deploy.service
[Unit]
Description=Run deploy script

[Service]
Type=oneshot
ExecStart=/home/ta1/src/my_website/scripts/deploy.sh
